# deploy/systemd/boxfleet-agent.service
[Unit]
Description=BoxFleet node agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStartPre=/opt/boxfleet/libexec/boxfleet-agent-guard guard --config /etc/boxfleet/agent.json
ExecStart=/opt/boxfleet/bin/boxfleet-agent run --config /etc/boxfleet/agent.json
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
