# internal/agent/templates/boxfleet-agent.service.tmpl
[Unit]
Description=BoxFleet node agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre={{ quote .AgentGuardPath }} guard --config {{ quote .AgentConfigPath }}
ExecStart={{ quote .AgentPath }} run --config {{ quote .AgentConfigPath }}
Restart={{ .Restart }}
RestartSec={{ .RestartSec }}

[Install]
WantedBy=multi-user.target
